# templates/zabbix-agent-supervisor.conf.erb
[program:zabbix-agent]
command=/usr/sbin/zabbix_agentd -c /etc/zabbix/zabbix_agentd.conf -f
autostart=true
autorestart=true
user=zabbix